# nebenkosten/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:postgresql]
command=/usr/lib/postgresql/16/bin/postgres -D /data/postgres -c listen_addresses=localhost -c unix_socket_directories=/var/run/postgresql
user=postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:backend]
command=/bin/bash -c "source /etc/nebenkosten.env && sleep 5 && cd /app/backend && alembic upgrade head && python scripts/seed.py && python scripts/reset_admin.py; uvicorn app.main:app --host 0.0.0.0 --port 8000"
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:frontend]
command=/bin/bash -c "source /etc/nebenkosten.env && node /app/frontend/server.js"
autostart=true
autorestart=true
priority=30
environment=PORT="3000",HOSTNAME="0.0.0.0",NODE_ENV="production"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
